# services/organizer.timer
[Unit]
Description=Timer for Organizer service

[Timer]
OnBootSec=2min
#OnUnitActiveSec=15sec

[Install]
WantedBy=timers.target

