# services/organizer.service
[Unit]
Description=Organizer service

[Service]
ExecStart=/usr/bin/organize run
Restart=always
RestartSec=20
User=evgeniy
Group=evgeniy

[Install]
WantedBy=multi-user.target

